# com.strix.telegram-bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.strix.telegram-bot</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/env</string>
        <string>python3</string>
        <string>-m</string>
        <string>strix_telegram_bot</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__BOT_DIR__</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>__BOT_DIR__/logs/strix_bot_stdout.log</string>

    <key>StandardErrorPath</key>
    <string>__BOT_DIR__/logs/strix_bot_stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>STRIX_BOT_DIR</key>
        <string>__BOT_DIR__</string>
        <key>STRIX_BIN</key>
        <string>__STRIX_BIN__</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
